# hack/docker/haproxy/1.8.4/sv/haproxy-controller/run
#!/bin/bash

source /etc/envvars

echo "Starting HAProxy controller ..."
cmd="exec voyager haproxy-controller $HAPROXY_CONTROLLER_ARGS"
echo $cmd
$cmd
