# hack/docker/haproxy/1.9.2/sv/socklog-unix/run
#!/bin/bash
exec 2>&1
exec chpst -Unobody socklog unix /dev/log
