# hack/docker/haproxy/1.8.3/sv/haproxy/run
#!/bin/bash
source /etc/envvars

exec haproxy -W -db -f /etc/haproxy/haproxy.cfg -p /var/run/haproxy.pid -sf $(cat /var/run/haproxy.pid)
